# root/etc/systemd/system/waptserver.service.d/nethserver.conf
[Unit]
Description=WAPT Server startup script custom for NethServer
After=network.target